# config/supervisor/laravel-queue-work.conf
[program:laravel-ai-orchestrator]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /mnt/linux-data/www/mira-park-manager/artisan queue:work --queue=default --sleep=3 --tries=3 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/mnt/linux-data/www/mira-park-manager/storage/logs/queue-worker.log
stopwaitsecs=3600
